# invest-bot.service
[Unit]
Description=Invest Bot - Telegram Investment Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/invest-bot
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/python3 /root/invest-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Настройки для логирования
SyslogIdentifier=invest-bot

[Install]
WantedBy=multi-user.target
