# invest-bot.service
[Unit]
Description=Investment Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/invest-bot
Environment=PATH=/root/invest-bot/venv/bin
ExecStart=/root/invest-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
